# monitoring-mixins/crontab
# applies any differences to your Grafana Mimir cluster
*/1 * * * * mimirtool rules sync --rule-dirs=/rules/